# scripts/supervisor/base.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:cloudera-scm-agent]
command=/bin/bash -c "exec /etc/init.d/cloudera-scm-agent start"
